# sys-apps/systemd-units/files/services-desktop/dhcpcd.service
[Unit]
Description=Dhcpcd

[Service]
ExecStart=/sbin/dhcpcd -q --nobackground

[Install]
WantedBy=network.target
