# sys-apps/systemd-service-files/files/services-desktop/dhcpcd.service
[Unit]
Description=Dhcpcd

[Service]
Type=forking
ExecStart=/sbin/dhcpcd -q
PIDFile=/var/run/dhcpcd.pid

[Install]
WantedBy=graphical.target
